# autostart/consciousness-services.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.s0fractal.consciousness</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chaoshex/s0fractal/autostart/start-all.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/chaoshex/s0fractal/logs/consciousness.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/chaoshex/s0fractal/logs/consciousness-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/chaoshex/.deno/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>